# recipes-tools/filesystem/files/chestnut/lib/systemd/system/chestnut-dig-fan.service
[Unit]
Description=Digital Board Fan Service

[Service]
ExecStart=/etc/chestnut/fpga_fan_script -s 120

[Install]
WantedBy=multi-user.target
