# recipes-tools/filesystem/files/cyan/lib/systemd/system/cyan-dig-fan.service
[Unit]
Description=Digital Board Fan Service

[Service]
ExecStart=/etc/cyan/fpga_fan_script -s 70

[Install]
WantedBy=multi-user.target
