# recipes-tools/crimson-filesystem/crimson-filesystem/lib/systemd/system/crimson-startup.service
[Unit]
Description=Crimson startup service
Before=getty@tty1.service
After=dev-tty1.device

[Service]
ExecStart=/etc/crimson/startup
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
